# ops/mac/shared/LaunchAgents/de.haraldweiss.ollama-tunnel-monitor.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>de.haraldweiss.ollama-tunnel-monitor</string>
    <key>ProgramArguments</key>
    <array>
        <string>/Users/haraldweiss/bin/check-ollama-tunnel.sh</string>
    </array>
    <key>StartInterval</key>
    <integer>300</integer>
    <key>RunAtLoad</key>
    <true/>
    <key>StandardOutPath</key>
    <string>/Users/haraldweiss/Library/Logs/ollama-tunnel-monitor.stdout.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/haraldweiss/Library/Logs/ollama-tunnel-monitor.stderr.log</string>
    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin:/opt/homebrew/bin</string>
    </dict>
</dict>
</plist>
