# ops/mac/mini/LaunchAgents/com.ai-provider.ollama-tunnel.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.ai-provider.ollama-tunnel</string>
    <key>ProgramArguments</key>
    <array>
        <string>/opt/homebrew/bin/autossh</string>
        <string>-M</string>
        <string>0</string>
        <string>-N</string>
        <string>-o</string>
        <string>ServerAliveInterval=30</string>
        <string>-o</string>
        <string>ServerAliveCountMax=3</string>
        <string>-o</string>
        <string>ExitOnForwardFailure=yes</string>
        <string>-o</string>
        <string>StrictHostKeyChecking=accept-new</string>
        <string>-R</string>
        <string>11435:127.0.0.1:11434</string>
        <string>-o</string>
        <string>IdentityFile=/Users/haraldweiss/.ssh/id_ed25519</string>
        <string>root@bewerbungen.wolfinisoftware.de</string>
    </array>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
    <key>StandardOutPath</key>
    <string>/Users/haraldweiss/Library/Logs/ai-provider-tunnel.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/haraldweiss/Library/Logs/ai-provider-tunnel.log</string>
    <key>EnvironmentVariables</key>
    <dict>
        <key>AUTOSSH_GATETIME</key>
        <string>30</string>
    </dict>
</dict>
</plist>
