# deploy/container/bewerbungen-imap-proxy.container
# SPDX-License-Identifier: AGPL-3.0-or-later
# © 2026 Harald Weiss
#
# Quadlet for the Bewerbungstracker IMAP Proxy (read-only).
# Binds only to host loopback — no external access.

[Unit]
Description=Bewerbungstracker IMAP Proxy (read-only)
After=network-online.target
Wants=network-online.target

[Container]
Image=localhost/bewerbungen:latest
ContainerName=bewerbungen-imap-proxy

Exec=imap-proxy

PublishPort=127.0.0.1:8765:8765

Volume=/opt/bewerbungen-data:/app/data:Z
Volume=/opt/bewerbungen-data/logs:/app/logs:Z

EnvironmentFile=/etc/bewerbungen/bewerbungen.env

Environment=FLASK_ENV=production
Environment=TZ=Europe/Berlin

Memory=256m
LogDriver=journald

[Service]
Restart=on-failure
RestartSec=10

[Install]
WantedBy=default.target
